# minecraft@.service
[Unit]
Description=Minecraft instance "%i"

[Service]
ExecStart=/usr/share/mc-server/run.sh %i
ExecStop=/usr/share/mc-server/stop.sh %i

[Install]
WantedBy=network.target
