# minecraft.target
[Unit]
Description=Minecraft Server
Requires=minecraft.pod
After=network-online.target

[Install]
WantedBy=multi-user.target