# minecraft.target
[Unit]
Description=Minecraft Server
Requires=minecraft-pod.service
After=network-online.target

[Install]
WantedBy=multi-user.target